# mqtt-chalet.service
[Unit]
Description=MQTT Chalet - DHT11 & RF Relay Controller
After=network.target mosquitto.service
Wants=mosquitto.service

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/MQTT-chalet
ExecStart=/home/pi/MQTT-chalet/venv/bin/python3 /home/pi/MQTT-chalet/main.py
Restart=on-failure
RestartSec=10
StandardOutput=append:/var/log/MQTT-chalet/output.log
StandardError=append:/var/log/MQTT-chalet/error.log
SyslogIdentifier=MQTT-chalet

[Install]
WantedBy=multi-user.target
